# meta-bsp/recipes-myir/board-info/board-info/myir-swupdate.service
[Unit]
Description=SWUpdate daemon
After=swupdate.service

[Service]
ExecStart=/etc/myir-swupdate.sh
KillMode=mixed

[Install]
WantedBy=multi-user.target
