# meta-bsp/recipes-myir/start-service/start-service/startup.service
[Unit]
Description=start process control system for UNIX
After=getty@tty1.service

[Service]
Type=forking
ExecStart=/etc/init.d/startup.sh
TimeoutSec=0

[Install]
WantedBy=multi-user.target
